# data/systemd/dcs-daq.service.in
[Unit]
Description=Data Acquisition Daemon
After=network.target
After=dbus.service
Conflicts=
Documentation=man:dcs-daq(1)

[Service]
EnvironmentFile=-@sysconfdir@/dcs
ExecStart=@sbindir@/dcs-daq $DAQD_ARGS
ExecReload=/bin/kill -HUP $MAINPID
# supress to log debug and error also to /var/log/messages
StandardOutput=null
StandardError=null
Type=dbus
BusName=org.opendcs.dcs.DAQ

[Install]
WantedBy=basic.target
Alias=dbus-org.opendcs.dcs.DAQ.service
